# debian/ec2-api-metadata.init.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          ec2-api-metadata
# Required-Start:    $network $local_fs $remote_fs $syslog
# Required-Stop:     $remote_fs
# Should-Start:      mysql postgresql rabbitmq-server keystone
# Should-Stop:       mysql postgresql rabbitmq-server keystone
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: EC2 Metadata API
# Description:       Provide Metadata API for ec2api
### END INIT INFO

# Author: Sam Morrison <sorrison@gmail.com>

PATH=/sbin:/usr/sbin:/bin:/usr/bin
DESC="OpenStack EC2API Metadata API"
PROJECT_NAME=ec2api
NAME=ec2-api-metadata
